# src/dist/feed-filter.service
[Unit]
Description=Feed Filter

[Service]
Type=simple
Restart=always
ExecStart=/opt/feedfilter/bin/feedfilter
KillMode=process
TimeoutSec=180
Environment=FEEDS=/opt/feedfilter/etc/feeds.yaml

[Install]
WantedBy=network.target